# supervisor/whatsapp-gateway.conf
[program:whatsapp-gateway]
command=node /var/www/html/whatsapp/index.js
directory=/var/www/html/whatsapp
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/whatsapp.log
stderr_logfile=/var/www/html/storage/logs/whatsapp_error.log
environment=NODE_ENV="production"

# Panduan Aktivasi Supervisor di VPS (Debian/Ubuntu):
# ----------------------------------------------------
# 1. Salin file konfigurasi ini ke direktori supervisor:
#    sudo cp supervisor/*.conf /etc/supervisor/conf.d/
#
# 2. Sesuaikan path "/var/www/html" dan "user=www-data" di dalam file .conf 
#    sesuai dengan direktori proyek dan user web server Anda di VPS.
#
# 3. Lakukan pembacaan ulang konfigurasi:
#    sudo supervisorctl reread
#
# 4. Tambahkan/aktifkan program baru:
#    sudo supervisorctl update
#
# 5. Cek status program:
#    sudo supervisorctl status
